# .external/etc/udev/rules.d/99-reload-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/su srdusr --command='systemctl --user start bspwm-reload.service'"
